# package_files/usr/share/doc/vscode-profile-importer/examples/vscode-profile-importer-sync.timer
[Unit]
Description=Run VSCode profile importer periodically (user)

[Timer]
OnCalendar=*-*-* *:00/6:00
Persistent=true

[Install]
WantedBy=timers.target
